# files/test2017.conf
description "uWSGI server instance configured to serve test2017"

start on runlevel [2345]
stop on runlevel [!2345]

setuid root
setgid www-data

env PATH=/home/ubuntu/venv/bin
chdir /home/ubuntu/test2017
exec uwsgi --ini uwsgi.ini
